# roles/base_hardening/templates/notify/on-boot.service.j2
# Rendered by Ansible — do not edit directly on the host.
[Unit]
Description=Notify on system boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ notify_script_dir }}/on-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
